# backend/template.config.ini
[db]
uri = mysql+pymysql://<username>:<password>@<host>/<dbname>

[webserver]
host = 127.0.0.1
port = 5000
secret_key = anything_you_wanna
log_path = /tmp
root_name = init_root_name
root_passwd = init_root_passwd


[uwsgi]
http = 0.0.0.0:5000
chdir = /path/to/backend
wsgi-file = app.py
callable = app
processes = 4
threads = 2

[program:blog]
command = env PATH="/path/to/your/python" uwsgi /path/to/config.ini
director = /path/to/backend
user = root
autostart = true
autorestart = true
stdout_logfile = /tmp/uwsgi_supervisor.log
die-on-term = true